# files/system/etc/profile.d/hachiman.sh
#! /usr/bin/env sh

# This sets the default location of the shared distrobox.ini file.
export GLOBAL_DISTROBOX_INI="/usr/share/hachiman/distrobox/distrobox.ini"
